# GENTOO-SYSTEMD/etc-systemd-gentoo-bak/local.d/apfs-mount.stop
#!/bin/bash


####/usr/bin/cpupower frequency-set -g performance
#/usr/sbin/x86_energy_perf_policy performance
#exec echo 0 > /sys/devices/system/cpu/intel_pstate/no_turbo
#/usr/bin/cpupower frequency-set -g performance
#mkdir -p /run/media/javashin/macos
#/usr/local/bin/apfs-fuse -o allow_other /dev/sda7 /run/media/javashin/macos

#fusermount 
/usr/bin/fusermount -u /run/media/javashin/macos
