# GENTOO-SYSTEMD/etc-systemd-gentoo-bak/local.d/libvirtd.start
#!/bin/bash
exec /bin/systemctl start libvirtd.service
