# GENTOO-SYSTEMD/etc-systemd-gentoo-bak/local.d/libvirtd.stop
#!/bin/bash
exec /bin/systemctl stop libvirtd.service
