# GENTOO-SYSTEMD/etc-systemd-gentoo-bak/local.d/turbo-boost.start
#!/bin/bash
/usr/bin/cpupower frequency-set -g schedutil
#performance
/usr/sbin/x86_energy_perf_policy performance
exec echo 1  > /sys/devices/system/cpu/cpufreq/boost
#exec echo 0 > /sys/devices/system/cpu/intel_pstate/no_turbo
#/usr/bin/cpupower frequency-set -g performance
