# GENTOO-SYSTEMD/etc-systemd-gentoo-bak/local.d/vmware.start
#!/bin/bash
exec /bin/systemctl start vmware-vmci.service vmware-vmmon.service vmware-vmnet.service vmware-vmsock.service vmware-usb.service

