# GENTOO-SYSTEMD/etc-systemd-gentoo-bak/local.d/zfs-arc.start
#!/bin/bash
echo 536870912 >> /sys/module/zfs/parameters/zfs_arc_max
/usr/local/bin/cleanram
# echo 0 >> /sys/module/zfs/parameters/zfs_arc_max
