# GENTOO-SYSTEMD/etc-systemd-gentoo-bak/profile.d/javashin-envs.sh
export PATH=/usr/lib/llvm/9/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/opt/bin

export LIBGL_DRI3_ENABLE=1
#export LIBGL_DEBUG=verbose
export mesa_glthread=true
export XDG_SESSION_TYPE=X11
export QT_QPA_PLATFORMTHEME=qt5ct
export XDG_SESSION_CLASS=X11
export EGL_PLATFORM=x11
export GDK_BACKEND=x11
export EDITOR=nano
