# GENTOO-SYSTEMD/etc-systemd-gentoo-bak/systemd/system/getty@tty1.service.d/override.conf
[Service]
Type=simple
ExecStart=
ExecStart=-/sbin/agetty --autologin javashin --noclear %I 38400 linux

