# GENTOO-SYSTEMD/etc-systemd-gentoo-bak/systemd/system/cpupower-frequency-set.service.d/00gentoo.conf
[Service]
Environment="OPTS=-g performance"
